# loki-addon/rootfs/etc/services.d/backup/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Backup Service..."

while true; do
    # Check for running backups using the Supervisor API instead of bashio::core.backup.is_running
    if bashio::api.hassio GET /backups/info | jq -e '.data.backups[] | select(.state == "running")' > /dev/null; then
        bashio::log.info "Home Assistant Backup is running, preparing Loki data..."
        
        # Temporary folder for backup
        mkdir -p /tmp/loki_backup
        
        # Copy important data
        cp -r /data/loki/wal /tmp/loki_backup/ || bashio::log.warning "Error copying WAL directory"
        cp -r /data/loki/chunks /tmp/loki_backup/ || bashio::log.warning "Error copying chunks directory"
        
        # Wait until backup is completed
        while bashio::api.hassio GET /backups/info | jq -e '.data.backups[] | select(.state == "running")' > /dev/null; do
            sleep 5
        done
        
        # Delete temporary backup files
        rm -rf /tmp/loki_backup
        
        bashio::log.info "Backup preparation completed"
    fi
    
    # Wait until next backup check
    sleep 3600 # Check hourly
done