# loki-addon/rootfs/etc/services.d/backup/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starte Backup-Service..."

while true; do
    # Prüfen, ob ein Home Assistant Backup läuft über die Supervisor API
    if bashio::api.supervisor GET /supervisor/info false | jq -e '.backups_in_progress > 0' > /dev/null; then
        bashio::log.info "Home Assistant Backup läuft, bereite Loki-Daten vor..."
        
        # Temporärer Ordner für Backup
        mkdir -p /tmp/loki_backup
        
        # Kopiere wichtige Daten
        cp -r /data/loki/wal /tmp/loki_backup/ || bashio::log.warning "Fehler beim Kopieren von wal"
        cp -r /data/loki/chunks /tmp/loki_backup/ || bashio::log.warning "Fehler beim Kopieren von chunks"
        
        # Warte bis das Backup abgeschlossen ist
        while bashio::api.supervisor GET /supervisor/info false | jq -e '.backups_in_progress > 0' > /dev/null; do
            sleep 5
        done
        
        # Lösche temporäre Backup-Dateien
        rm -rf /tmp/loki_backup
        
        bashio::log.info "Backup-Vorbereitung abgeschlossen"
    fi
    
    # Warte bis zum nächsten Backup-Zeitpunkt
    sleep 3600 # Prüfe stündlich
done