# loki-addon/rootfs/etc/services.d/backup/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starte Backup-Service..."

while true; do
    # Prüfen, ob ein Home Assistant Backup läuft
    if bashio::supervisor.ping; then
        # Verwende API, um Backup-Status zu prüfen
        BACKUP_RUNNING=$(bashio::api.supervisor GET /backups false ".backups[] | select(.state == \"running\")" | wc -l)
        
        if [[ $BACKUP_RUNNING -gt 0 ]]; then
            bashio::log.info "Home Assistant Backup läuft, bereite Loki-Daten vor..."
            
            # Temporärer Ordner für Backup
            mkdir -p /tmp/loki_backup
            
            # Kopiere wichtige Daten
            cp -r /data/loki/wal /tmp/loki_backup/ || true
            cp -r /data/loki/chunks /tmp/loki_backup/ || true
            
            # Warte bis keine laufenden Backups mehr vorhanden sind
            while true; do
                BACKUP_RUNNING=$(bashio::api.supervisor GET /backups false ".backups[] | select(.state == \"running\")" | wc -l)
                if [[ $BACKUP_RUNNING -eq 0 ]]; then
                    break
                fi
                sleep 5
            done
            
            # Lösche temporäre Backup-Dateien
            rm -rf /tmp/loki_backup
            
            bashio::log.info "Backup-Vorbereitung abgeschlossen"
        fi
    fi
    
    # Warte bis zum nächsten Backup-Zeitpunkt
    sleep 3600 # Prüfe stündlich
done